# scripts/com.klab.folder-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.klab.folder-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/backup-dev-apps.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardErrorPath</key>
    <string>__HOME__/.local/logs/backup-dev-apps-stderr.log</string>
</dict>
</plist>
